# supervisor/generate_blowfish_secret.conf
[program:generate_blowfish_secret]
directory=/opt/phpmyadmin/src
command=bash generate_blowfish_secret.sh
startsecs=0
