# supervisor/tunnel.conf
[program:tunnel]
directory=/opt/phpmyadmin/src/
command=python3 tunnel.py
startretries=999
startsecs=0
