# supervisor/use_custom_config.conf
[program:use_custom_config]
directory=/opt/phpmyadmin/src
command=php -f use_custom_config.php
startsecs=0
